# containers/systemd/caddy.container
[Unit]
Description=Caddy — reverse proxy, static SPA host, OHIF viewer origin (spec §2.1, §3, §9)
After=network-online.target

[Container]
Image=docker.io/library/caddy:2
ContainerName=caddy
Network=meld-net.network
Volume=caddy-data.volume:/data
Volume=%h/.config/meld7t/caddy/Caddyfile:/etc/caddy/Caddyfile:ro,z
# Built SPA shell bundle (`just web-build` → platform/web/dist); Caddy serves it at / (§9.4).
Volume=%h/meld7t/platform/web/dist:/srv/spa:ro,z
# Bound to all interfaces (0.0.0.0) so the demo is reachable across the intranet/home LAN.
# 9443 = shell + API; 9444 = OHIF viewer origin. HTTP catch-all (see Caddyfile; HTTPS is the
# documented production hardening). Open these ports in the host firewall for remote devices.
PublishPort=9443:9443
PublishPort=9444:9444
NoNewPrivileges=true

[Service]
Restart=always

[Install]
WantedBy=default.target
